# deployment/ai-news-crawler.service
[Unit]
Description=AI News Crawler - Daily university AI news monitoring with auto-commit
Documentation=https://github.com/tyson-swetnam/webcrawler
After=network.target

[Service]
Type=oneshot
User=tswetnam
Group=tswetnam
WorkingDirectory=/home/tswetnam/github/webcrawler

# Environment
Environment="PATH=/home/tswetnam/github/webcrawler/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/tswetnam/github/webcrawler/.env

# Execution - runs crawler and commits/pushes changes
ExecStart=/home/tswetnam/github/webcrawler/scripts/run_crawler_and_commit.sh

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-news-crawler

# Resource limits
TimeoutStartSec=1800
TimeoutStopSec=30

# Restart policy
Restart=no

[Install]
WantedBy=multi-user.target
