# deployment/ai-news-crawler.service
[Unit]
Description=AI News Crawler - Daily university AI news monitoring
Documentation=https://github.com/yourorg/ai-news-crawler
After=network.target postgresql.service

[Service]
Type=oneshot
User=crawler
Group=crawler
WorkingDirectory=/opt/ai-news-crawler

# Environment
Environment="PATH=/opt/ai-news-crawler/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/ai-news-crawler/.env

# Execution
ExecStart=/opt/ai-news-crawler/venv/bin/python -m crawler

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-news-crawler

# Resource limits
TimeoutStartSec=600
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ai-news-crawler /var/log/ai-news-crawler

# Restart policy
Restart=no

[Install]
WantedBy=multi-user.target
